# ops/overlay-routing-healer/overlay-routing-healer.timer
[Unit]
Description=Run overlay routing health check every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
