# ops/mailbox-watch-renewer/mailbox-watch-renewer.service
[Unit]
Description=Renew Gmail/Microsoft mailbox watch subscriptions (email-ingress)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mailbox-watch-renewer.sh
